# Downloads/idr-inostudio-final/src/idr-src/idr-config/servers/idr-api.staging/files/etc/cron.daily/renew-ssl
#!/bin/bash
MAILTO=""
/usr/bin/certbot-auto certonly -c /www/letsencrypt/idr-api.ini --renew-by-default
sleep 5; service nginx restart
