# Downloads/idr-inostudio-final/src/idr-src/idr-config/servers/idr.staging/files/etc/cron.daily/renew-ssl
#!/bin/bash
/usr/bin/certbot-auto certonly -c /www/letsencrypt/idr.ini --renew-by-default
sleep 5; service nginx restart
